# server/archiver/install
cp collector.service /etc/systemd/system/archiver-collector.service
systemctl enable archiver-collector
systemctl start archiver-collector
cp collector_snapshots.crontab /etc/cron.d/collector-snapshots
/etc/init.d/cron reload
