# server/config/cron.daily
#!/bin/sh

if [ -f /opt/weather/scripts/daily ]; then
        /bin/sh /opt/weather/scripts/daily &
fi

exit 0

