# server/config/mainview-server-http.service
[Unit]
Description=Weather Mainview Server HTTP Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/weather/server/http/server-mainview.js /opt/weather/server/secrets.txt
ExecReload=kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

