# server/config/weather-daily.timer
[Unit]
Description=Run Weather Daily Script

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
