# archiver/install
cp collector.service /etc/systemd/system/archiver_collector.service
systemctl enable archiver_collector
systemctl start archiver_collector
cp collector_snapshots.crontab /etc/cron.d/collector_snapshots
/etc/init.d/cron reload
